# supervisord.conf
[supervisord]
nodaemon=true
autorestart=false
[program:userdetails]
command=python userdetails.py
[program:api]
command=python api.py -DFOREGROUND

